# config/bearings.service
[Unit]
Description=Bearings — Claude Code agent web UI
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/bearings serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
